# linear-agent.service
[Unit]
Description=Linear AI Agent — Hermes-powered autonomous agent for Linear
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/abe/linear-agent
EnvironmentFile=/home/abe/linear-agent/.env
ExecStartPre=/bin/sh -c 'while ss -tlnp sport = :8660 2>/dev/null | grep -q :8660; do sleep 1; done'
ExecStart=/home/abe/linear-agent/.venv/bin/uvicorn linear_agent:app --host 0.0.0.0 --port 8660 --log-level info
Restart=always
RestartSec=15
TimeoutStopSec=15

# Security hardening (user service compatible)
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/home/abe/linear-agent/workspace

[Install]
WantedBy=multi-user.target
