# linear-agent-user.service
[Unit]
Description=Linear AI Agent — Hermes-powered autonomous agent for Linear
After=network.target network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/home/abe/linear-agent
EnvironmentFile=/home/abe/linear-agent/.env
ExecStartPre=/bin/sh -c 'while ss -tlnp sport = :8660 2>/dev/null | grep -q :8660; do sleep 1; done'
ExecStart=/home/abe/linear-agent/bin/linear-agent-wrapper.sh
Restart=always
RestartSec=15
TimeoutStopSec=15

# Security hardening (user service compatible)
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/home/abe/linear-agent/workspace

# systemd-notify based health checks (wrapper uses systemd-notify --ready + WATCHDOG=1)
Type=notify
WatchdogSec=30
NotifyAccess=main

[Install]
WantedBy=multi-user.target
